# crontab
# Twilio phone test — every 10 min during working hours (Fortaleza/CE, BRT)
# Runs at :00 :10 :20 :30 :40 :50 of each hour from 08:00 to 18:50 Mon–Fri
*/10 8-18 * * 1-5 root /usr/local/bin/python /app/test_call.py >> /var/log/twilio-test.log 2>&1
